# indicateurs/indicateurs.timer
[Unit]
Description="timer pour la maj des indicateurs"
Requires=indicateurs.service
[Timer]
Unit=indicateurs.service
OnCalendar=*-*-* *:0/15:00
AccuracySec=20s

[Install]
WantedBy=timers.target